# packages/cli/templates/secrets-portal.service.j2
[Unit]
Description=HostKit Secrets Portal
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/hostkit/secrets-portal
Environment=MASTER_KEY_PATH=/etc/hostkit/master.key
Environment=SECRETS_DIR=/var/lib/hostkit/secrets
Environment=HOSTKIT_DB_PATH=/var/lib/hostkit/hostkit.db
Environment=PORTAL_HOST=127.0.0.1
Environment=PORTAL_PORT=8900
ExecStart=/var/lib/hostkit/secrets-portal/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8900
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
